# recipes-flir/videorender/files/videorender.service
[Unit]
Description=Videorender application
After=weston.service

[Service]
Type=simple
ExecStart=/FLIR/usr/bin/videorender.sh
Environment=LD_LIBRARY_PATH=/FLIR/usr/lib

[Install]
WantedBy=multi-user.target
