# recipes-flir/sysmon/files/sysmon.timer
[Unit]
Description=FLIR Watchdog service start timer

[Timer]
OnBootSec=1min
Unit=sysmon.service

[Install]
WantedBy=multi-user.target

