# recipes-flir/skylab/files/skylab-cmd.service
[Unit]
Description=Skylab cmd
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/skylab-cmd \
			      --verbosity=5 \
			      --server-connect-endpoint=tcp://127.0.0.1:5552 \
			      --client-bind-endpoint=tcp://*:5551
Restart=on-failure

[Install]
WantedBy=multi-user.target
