# recipes-flir/skylab/files/skylab-streamserver.service
[Unit]
Description=Skylab streamserver
After=network.target bluetooth.target

[Service]
Type=dbus
BusName=com.flir.skylab1.streamserver
ExecStart=/usr/bin/skylab-streamserver
Restart=on-failure
Environment=LD_LIBRARY_PATH=/FLIR/usr/lib
ExecStopPost=/sbin/systemd-failure-notifier "%n" "%i"

[Install]
WantedBy=multi-user.target
