# recipes-connectivity/rndisload/files/rndis.service
[Unit]
Description=Initialize RNDIS on usb0

[Service]
Type=oneshot
ExecStart=/sbin/rndisload
ExecStop=/sbin/rndisunload
RemainAfterExit=yes
