# recipes-flir/flir-activity-monitor/files/flir-activity-monitor.service
[Unit]
Description=FLIR Activity Monitor
After=flirapp.service

[Service]
Type=simple
ExecStart=/usr/bin/flir-activity-monitor

[Install]
WantedBy=multi-user.target

